# nams-e2-micro/chromium/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1024x640x24 -ac +extension RANDR
priority=10
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err.log

[program:fluxbox]
command=/usr/bin/fluxbox -display :99
priority=20
autostart=true
autorestart=true
environment=DISPLAY=":99",HOME="/data"
stdout_logfile=/var/log/supervisor/fluxbox.log
stderr_logfile=/var/log/supervisor/fluxbox.err.log

[program:chromium]
command=/usr/bin/chromium --no-sandbox --disable-dev-shm-usage --disable-gpu --no-first-run --no-default-browser-check --disable-background-networking --disable-background-timer-throttling --disable-component-update --disable-default-apps --disable-extensions --disable-sync --disable-features=Translate,BackForwardCache,AcceptCHFrame,MediaRouter --renderer-process-limit=1 --js-flags=--max-old-space-size=192 --remote-debugging-address=127.0.0.1 --remote-debugging-port=9222 --user-data-dir=/data/profile --window-size=1024,640 about:blank
priority=30
autostart=true
autorestart=true
environment=DISPLAY=":99",HOME="/data"
stdout_logfile=/var/log/supervisor/chromium.log
stderr_logfile=/var/log/supervisor/chromium.err.log

[program:cdp-proxy]
command=/usr/bin/socat TCP-LISTEN:9223,bind=0.0.0.0,reuseaddr,fork TCP:127.0.0.1:9222
priority=35
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/var/log/supervisor/cdp-proxy.log
stderr_logfile=/var/log/supervisor/cdp-proxy.err.log

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -xkb
priority=40
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err.log

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 6080 localhost:5900
priority=50
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err.log
